# systemd/dockerhook.service
[Unit]
Description=DockerHook front-end (express)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/DockerHook
Environment=NODE_ENV=production
Environment=TRIGGER_DIR=/var/run/dockerhook
ExecStart=/usr/bin/node /root/DockerHook/src/index.js
Restart=on-failure

[Install]
WantedBy=multi-user.target
