# systemd/watcher.service
[Unit]
Description=DockerHook trigger watcher
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /root/dockerhook/src/watcher.js
Restart=on-failure

[Install]
WantedBy=multi-user.target
