# systemd/watcher.service
[Unit]
Description=DockerHook trigger watcher
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/DockerHook
Environment=NODE_ENV=production
Environment=TRIGGER_DIR=/var/run/dockerhook
Environment=COMPOSE_DIR=/root/
ExecStart=/usr/bin/node /root/DockerHook/src/watcher.js
Restart=on-failure

[Install]
WantedBy=multi-user.target
